# Build/nautilus-data.service
[Unit]
Description=Nautilus algorithmic trading data service
After=network.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/NautilusDataRelease/
ExecStart=/usr/bin/dotnet /home/ec2-user/NautilusDataRelease/NautilusData.dll
Restart=always
RestartSec=10
SyslogIdentifier=nautilus-data

[Install]
WantedBy=multi-user.target
